# scripts/skel/debian/new_service.timer.j2
[Unit]
Description = My Service

[Timer]
OnBootSec = 15min
OnUnitActiveSec = 1d
Persistent=true

[Install]
WantedBy = timers.target
